# systemd/gambota-sonoff-server.service
[Unit]
Description=Service Gambota Sonoff Server

[Service]
ExecStart=/opt/gambota-sonoff-server/webserver.py

[Install]
WantedBy=multi-user.target
